# launchagents/com.klaudioz.openportal-dashboard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.klaudioz.openportal-dashboard</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/opt/node/bin/node</string>
        <string>/Users/klaudioz/.local/share/openportal/dashboard.js</string>
    </array>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOSTNAME</key>
        <string>m4-mini.tail09133d.ts.net</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/klaudioz/.local/share/openportal/dashboard.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/klaudioz/.local/share/openportal/dashboard.log</string>
</dict>
</plist>
